# packaging/root/DEBIAN/postinst
#!/bin/sh
set -e
echo 'Installing python3 dependencies...'
pip3 install -r /usr/share/stoneflashtool/requirements.txt
echo 'Setting up service...'
systemctl enable --now stoneflashtool.service
